# services/notification-service/deployment/eracun-notification-service.service
[Unit]
Description=eRacun Notification Service
Documentation=https://github.com/eracun/notification-service
After=network-online.target rabbitmq-server.service postgresql.service
Wants=network-online.target
PartOf=eracun.target

[Service]
Type=simple
User=eracun
Group=eracun

# Working directory
WorkingDirectory=/opt/eracun/services/notification-service

# Environment configuration
EnvironmentFile=/etc/eracun/notification-service.env

# Service execution
ExecStart=/usr/bin/node dist/index.js

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Resource limits
MemoryMax=2G
MemoryHigh=1.5G
CPUQuota=200%

# Security hardening
# Filesystem protection
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log/eracun

# No new privileges
NoNewPrivileges=true

# Restrict namespaces
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount
SystemCallErrorNumber=EPERM

# Network restrictions (none - service needs network access for SMTP/Twilio)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Restrict realtime
RestrictRealtime=true

# Hide sensitive paths
InaccessiblePaths=/etc/eracun/.age-key

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eracun-notification-service

# Process properties
Nice=0
OOMScoreAdjust=-100

[Install]
WantedBy=eracun.target multi-user.target
