# deployment/systemd/eracun-fina-connector.service
[Unit]
Description=eRacun FINA Connector Service
Documentation=https://github.com/tpatarci/eRacun-development
After=network-online.target postgresql.service rabbitmq-server.service eracun-digital-signature-service.service
Wants=network-online.target
Requires=postgresql.service rabbitmq-server.service

[Service]
Type=simple
User=eracun
Group=eracun
WorkingDirectory=/opt/eracun/services/fina-connector

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/etc/eracun/fina-connector.env

# Execution
ExecStartPre=/usr/local/bin/sops-decrypt.sh fina-connector
ExecStart=/usr/bin/node dist/index.js
ExecStop=/bin/kill -SIGTERM $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=300s
StartLimitBurst=5

# Resource limits
MemoryMax=1G
MemoryHigh=768M
CPUQuota=200%
TasksMax=512

# Security Hardening
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/eracun/fina-connector /var/log/eracun
ReadOnlyPaths=/etc/eracun/secrets/certs

NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
SecureBits=keep-caps-locked noroot noroot-locked no-setuid-fixup no-setuid-fixup-locked

SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete @debug @mount @swap @reboot @module @raw-io
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native

RestrictAddressFamilies=AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=10.0.0.0/8
IPAddressAllow=172.16.0.0/12
IPAddressAllow=192.168.0.0/16
# Allow FINA test/prod endpoints
IPAddressAllow=195.29.162.0/24

ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true

LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RemoveIPC=true
PrivateMounts=true

InaccessiblePaths=/etc/eracun/.age-key
UMask=0077

StandardOutput=journal
StandardError=journal
SyslogIdentifier=eracun-fina-connector

[Install]
WantedBy=multi-user.target
