# deployment/systemd/eracun-digital-signature-service.service
# eRacun Digital Signature Service
#
# XMLDSig enveloped signature generation and verification for invoices.
# High-throughput batch signing with RSA-SHA256 algorithm.
#
# INSTALLATION:
# sudo cp eracun-digital-signature-service.service /etc/systemd/system/
# sudo systemctl daemon-reload
# sudo systemctl enable eracun-digital-signature-service
# sudo systemctl start eracun-digital-signature-service

[Unit]
Description=eRacun Digital Signature Service (XMLDSig)
Documentation=https://github.com/eracun/eRacun-development
After=network-online.target rabbitmq-server.service eracun-cert-lifecycle-manager.service
Wants=network-online.target rabbitmq-server.service

[Service]
Type=simple
User=eracun
Group=eracun
WorkingDirectory=/opt/eracun/services/digital-signature-service

# Environment configuration
EnvironmentFile=/etc/eracun/platform.conf
EnvironmentFile=/etc/eracun/environment.conf
EnvironmentFile=/etc/eracun/services/digital-signature-service.conf
EnvironmentFile=/run/eracun/secrets.env

# Decrypt secrets before starting
ExecStartPre=/usr/local/bin/decrypt-secrets.sh digital-signature-service

# Start service
ExecStart=/usr/local/bin/node dist/index.js

# Restart policy (critical service)
Restart=on-failure
RestartSec=5s
StartLimitBurst=10
StartLimitIntervalSec=300s

# Resource limits (higher CPU for crypto operations)
MemoryMax=2G
MemoryHigh=1600M
CPUQuota=400%
TasksMax=1024

# ============================================================================
# Systemd Hardening
# ============================================================================

# Filesystem Protection
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log/eracun/digital-signature-service /var/lib/eracun/certs
ReadOnlyPaths=/etc/eracun/secrets

# Prevent Privilege Escalation
NoNewPrivileges=true

# Prevent Access to Sensitive Files
InaccessiblePaths=/etc/eracun/.age-key /root /home

# Namespace Isolation
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true

# Capability Restrictions
CapabilityBoundingSet=
AmbientCapabilities=

# System Call Filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete @debug
SystemCallErrorNumber=EPERM

# Restrict Address Families
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Restrict Namespaces
RestrictNamespaces=true

# Lock Personality
LockPersonality=true

# Memory Deny Write Execute (disable for crypto operations if needed)
# MemoryDenyWriteExecute=true

# Restrict Realtime
RestrictRealtime=true

# Remove IPC
RemoveIPC=true

# Network Access (internal only)
# IPAddressDeny=any
# IPAddressAllow=localhost 10.0.0.0/8

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eracun-digital-signature-service

[Install]
WantedBy=multi-user.target
