# services/kpd-registry-sync/kpd-registry-sync.service
[Unit]
Description=KPD Registry Sync Service - KLASUS product classification code sync
Documentation=https://github.com/eracun/eracun-platform
After=network-online.target postgresql.service rabbitmq.service
Wants=network-online.target

[Service]
Type=simple
User=eracun
Group=eracun
WorkingDirectory=/opt/eracun/services/kpd-registry-sync

# Environment
EnvironmentFile=/etc/eracun/kpd-registry-sync.env

# Execution
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=600

# Resource limits
MemoryLimit=512M
MemoryHigh=400M
CPUQuota=100%
TasksMax=512

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
CapabilityBoundingSet=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native

# Writable directories (logs, data)
ReadWritePaths=/var/log/eracun/kpd-registry-sync
ReadWritePaths=/var/lib/eracun/kpd-registry-sync

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kpd-registry-sync

# Graceful shutdown (30 second timeout)
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
