# services/xsd-validator/xsd-validator.service
[Unit]
Description=eRacun XSD Validator Service
Documentation=https://github.com/eracun/eracun-development
After=network.target rabbitmq-server.service
Wants=rabbitmq-server.service

[Service]
Type=simple
User=eracun
Group=eracun
WorkingDirectory=/opt/eracun/services/xsd-validator

# Environment variables
Environment="NODE_ENV=production"
Environment="LOG_LEVEL=info"
Environment="RABBITMQ_URL=amqp://localhost:5672"
Environment="RABBITMQ_QUEUE=eracun.xsd-validator.xml"
Environment="SCHEMA_PATH=/opt/eracun/services/xsd-validator/schemas/ubl-2.1"
Environment="PROMETHEUS_PORT=9100"
Environment="HEALTH_PORT=8080"
Environment="JAEGER_ENDPOINT=http://localhost:14268/api/traces"
Environment="MAX_CONCURRENT=100"

# Start command
ExecStart=/usr/bin/node dist/index.js

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60s

# Resource limits (per TODO-008 and README)
MemoryMax=256M
MemoryHigh=200M
CPUQuota=100%

# Security hardening (per TODO-008)
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Read-only paths
ReadOnlyPaths=/opt/eracun/services/xsd-validator/schemas

# Writable paths (logs, temp files)
ReadWritePaths=/var/log/eracun
StateDirectory=eracun/xsd-validator
CacheDirectory=eracun/xsd-validator

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete
SystemCallErrorNumber=EPERM

# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# Networking
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=10.0.0.0/8
IPAddressAllow=172.16.0.0/12
IPAddressAllow=192.168.0.0/16

# File descriptor limit
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
